# debian/postrm
if [ "$1" = "purge" ] ; then
	update-rc.d hid-mapper remove >/dev/null
fi
